# images/cribl/container-files/etc/s6-overlay/s6-rc.d/cribl/run
#!/command/execlineb

with-contenv

# Tell s6-supervise to only permit one running instance
# foreground { s6-svc -O /run/service/cribl }

# s6-setuidgid nonroot

redirfd -w -n 1 /dev/null
redirfd -w -n 2 /dev/null

exec wait-all /opt/cribl/bin/cribl server
